# docker/hyperbricks-deploy
#!/sbin/openrc-run

name="hyperbricks-deploy"
description="Hyperbricks Deploy API"
command="/usr/local/bin/hyperbricks"
command_args="start --deploy-remote"
command_user="deploy:deploy"
pidfile="/run/${name}.pid"
directory="/opt/hyperbricks"
command_background="yes"

depend() {
  need net
}
